# systemd/first-boot-wizard.service
[Unit]
Description=Jarvis First Boot Wizard
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/jarvis-first-boot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
