# systemd/jarvis-greeting.service
[Unit]
Description=J.A.R.V.I.S. Startup Greeting
Documentation=https://github.com/jarvis/jarvis
# Run after the graphical session and network are available
After=graphical-session.target network-online.target
Wants=graphical-session.target
# Don't fail the session if this fails
FailureAction=none

[Service]
Type=oneshot
RemainAfterExit=no
# Inherit the user's display/audio environment
Environment=DISPLAY=:0
Environment=PULSE_SERVER=unix:/run/user/%U/pulse/native
EnvironmentFile=-/etc/jarvis/config.env
ExecStart=/opt/jarvis/scripts/startup_greeting.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=60

[Install]
WantedBy=default.target
