# systemd/jarvis.service
[Unit]
Description=Jarvis Assistant Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/jarvis
EnvironmentFile=/etc/jarvis/config.env
ExecStart=/opt/jarvis/scripts/run_jarvis.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
